# dotfiles/other-linux/.config/systemd/user/wm.target
[Unit]
Description=Window Manager
AllowIsolate=yes
