# dotfiles/config/systemd/user/compton.service
[Unit]
Description=compton

[Service]
Type=simple
ExecStart=/usr/bin/compton -cC
ExecStop=/usr/bin/pkill compton

[Install]
WantedBy=wm.target
