# dotfiles/config/systemd/user/copyq.service
[Unit]
Description=CopyQ

[Service]
Type=simple
ExecStart=/usr/bin/copyq
ExecStop=/usr/bin/pkill copyq

[Install]
WantedBy=wm.target
