# dotfiles/config/systemd/user/dunst.service
[Unit]
Description=Dunst
After=taffybar.service

[Service]
Type=simple
ExecStart=/usr/bin/dunst
ExecStop=/usr/bin/pkill dunst

[Install]
WantedBy=wm.target
