# dotfiles/config/systemd/user/mopidy.service
[Unit]
Description=mopidy

[Service]
Type=simple
ExecStart=/usr/bin/env mopidy
ExecStop=/usr/bin/env pkill mopidy

[Install]
WantedBy=wm.target
