# dotfiles/config/systemd/user/parcellite.service
[Unit]
Description=Parcellite clipboard manager

[Service]
ExecStart=/usr/bin/parcellite
Restart=on-abort

[Install]
WantedBy=graphical.target
