# dotfiles/config/systemd/user/pasystray.service
[Unit]
Description=Pulse Audio System Tray

[Service]
Type=simple
ExecStart=/usr/bin/pasystray
ExecStop=/usr/bin/pkill pasystray

[Install]
WantedBy=wm.target
