# dotfiles/config/systemd/user/picom.service
[Unit]
Description=picom

[Service]
Type=simple
ExecStart=/usr/bin/env picom -cC
ExecStop=/usr/bin/env pkill picom

[Install]
WantedBy=wm.target
