# dotfiles/config/systemd/user/taffybar.service
[Unit]
Description=taffybar

[Service]
Type=simple
ExecStart=/usr/bin/env taffybar
ExecStop=/usr/bin/pkill taffybar

[Install]
WantedBy=wm.target
