# dotfiles/config/systemd/user/udiskie.service
[Unit]
Description=udiskie
After=taffybar.service

[Service]
ExecStart=/usr/bin/udiskie --tray
Restart=always

[Install]
WantedBy=wm.target
