# dotfiles/config/systemd/user/udiskie.service
[Unit]
Description=udiskie

[Service]
ExecStart=/usr/bin/udiskie --tray
Restart=always

[Install]
WantedBy=default.target
