# dotfiles/config/systemd/user/volnoti.service
[Unit]
Description=volnoti

[Service]
Type=forking
ExecStart=/usr/bin/volnoti -nw
ExecStop=/usr/bin/pkill volnoti

[Install]
WantedBy=wm.target
