# dotfiles/config/systemd/user/wallpaper.service
[Unit]
Description=wallpaper

[Service]
Type=oneshot
ExecStart=/usr/bin/env wallpaper.sh

[Install]
WantedBy=wm.target
