# resources/LaunchAgents/org.imalison.mu4e-alerts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.imalison.mu4e-alerts</string>
    <key>ProgramArguments</key>
    <array>
      <string>zsh</string>
      <string>-c</string>
        <string>refresh_config && mu4e_update_index_and_alert</string>
    </array>
    <key>StartInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/tmp/org.imalison.mu4e-alerts.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/org.imalison.mu4e-alerts.error</string>
</dict>
</plist>
